# scripts/com.cfte.live-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cfte.live-loop</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nguyenduythanh/Downloads/crypto-flow-thesis-engine-starter/crypto-flow-thesis-engine-starter/.venv/bin/python3</string>
        <string>-m</string>
        <string>cfte.cli.main</string>
        <string>run-live</string>
        <string>--symbol</string>
        <string>BTCUSDT</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/nguyenduythanh/Downloads/crypto-flow-thesis-engine-starter/crypto-flow-thesis-engine-starter</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>src</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/tmp/cfte_live_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cfte_live_stderr.log</string>
</dict>
</plist>
